# agent1/deployment/agent1-production.service
[Unit]
Description=OpenPredict Agent 1 Data Collector (Production)
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/openpredict-prod/agent1
Environment="PATH=/home/ubuntu/openpredict-prod/agent1/venv/bin"
EnvironmentFile=/home/ubuntu/openpredict-prod/agent1/.env.production
ExecStart=/home/ubuntu/openpredict-prod/agent1/venv/bin/python agent1_main.py
Restart=on-failure
RestartSec=10
StandardOutput=append:/var/log/agent1/agent1_production.log
StandardError=append:/var/log/agent1/agent1_production_error.log

# Resource limits
MemoryLimit=1G
CPUQuota=50%

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
